# etc/systemd/system/zerotier-one.service.d/override.conf
[Service]
    # Remove any default ExecStop command
    ExecStop=
    # This is the registration: Run the script *before* the service starts
    ExecStartPre=/usr/local/bin/zerotier-cleanup
